# ansible/templates/protosuit-networkingbridge.service.j2
[Unit]
Description=Protosuit Engine Networking Bridge
After=network.target mosquitto.service NetworkManager.service
Wants=mosquitto.service NetworkManager.service
Requires=mosquitto.service

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ project_dir }}
Environment=PYTHONPATH={{ project_dir }}
Environment=PYTHONUNBUFFERED=1
ExecStart={{ project_dir }}/env/bin/python {{ project_dir }}/networkingbridge/networkingbridge.py
Restart={{ service_config.restart_policy }}
RestartSec={{ service_config.restart_sec }}
TimeoutStopSec={{ service_config.timeout_stop_sec }}
StandardOutput=journal
StandardError=journal

# Security settings
# NoNewPrivileges must be false to allow sudo for network management
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=read-only
ReadWritePaths={{ project_dir }}
ReadWritePaths=/home/{{ service_user }}/.config
ReadWritePaths=/home/{{ service_user }}/.cache
ReadWritePaths=/tmp

# Allow access to network management
SupplementaryGroups=netdev

[Install]
WantedBy=graphical.target
