# ansible/templates/protosuit-controllerbridge.service.j2
[Unit]
Description=Protosuit Engine Controller Bridge (Gamepad Input)
After=network.target mosquitto.service
Wants=mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ project_dir }}
Environment=PYTHONPATH={{ project_dir }}/engine
Environment=PYTHONUNBUFFERED=1
ExecStart={{ project_dir }}/env/bin/python {{ project_dir }}/engine/controllerbridge/controllerbridge.py
Restart={{ service_config.restart_policy }}
RestartSec={{ service_config.restart_sec }}
TimeoutStopSec={{ service_config.timeout_stop_sec }}
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ project_dir }}
ReadWritePaths=/home/{{ service_user }}/.config
ReadWritePaths=/home/{{ service_user }}/.cache
ReadWritePaths=/tmp
ReadWritePaths=/dev/input

# Allow access to input devices
SupplementaryGroups=input

[Install]
WantedBy=graphical.target
