# ansible/templates/protosuit-systembridge.service.j2
[Unit]
Description=Protosuit Engine System Bridge (Metrics, Fan, Power)
After=network.target mosquitto.service
Wants=mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ project_dir }}
Environment=PYTHONPATH={{ project_dir }}/engine
Environment=PYTHONUNBUFFERED=1
ExecStart={{ project_dir }}/env/bin/python {{ project_dir }}/engine/systembridge/systembridge.py
Restart={{ service_config.restart_policy }}
RestartSec={{ service_config.restart_sec }}
TimeoutStopSec={{ service_config.timeout_stop_sec }}
StandardOutput=journal
StandardError=journal

# Security settings
# NoNewPrivileges=false because sudo tee is used for sysfs/config writes
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths={{ project_dir }}
ReadWritePaths=/tmp

# Allow access to vcgencmd (video group)
SupplementaryGroups=video

[Install]
WantedBy=graphical.target
