# ansible/templates/xserver.service.j2
[Unit]
Description=X Server for Protosuit Engine
After=systemd-modules-load.service sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/start-x-server.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=graphical.target
