# ansible/roles/hub_networking_setup/templates/usb-interface-monitor.service.j2
[Unit]
Description=USB Interface Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/check-usb-interfaces.sh
User=root
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target